# DOCKER/X64/docker_memcached/assets/memcached_supervisor.conf
[supervisord]
nodaemon=true

[program:memcached]
command=memcached  -vv -m 256M -I 1M -f 3 -p 11211 -u memcache -l 0.0.0.0
